# deb/vmware-horizonagent-linux-x86_64-7.0.3/opt/vmware/VMware-horizonagent-linux-x86_64-7.0.3-4679208-config/upgradeHorizonAgent.service
[Unit]
Description=upgrade VMWare Horizon Agent 
After=network-online.target
After=multi-user.target

[Service]
Type=idle
User=root
PIDFile=/run/upgradeHorizonAgent.pid
ExecStart=/usr/sbin/install_viewagent -A yes -M yes

[Install]
WantedBy=multi-user.target