# deb/vmware-horizonagent-linux-x86_64-7.1.0/opt/vmware/VMware-horizonagent-linux-x86_64-7.1.0-5141896-config/upgradeHorizonAgent.conf
description "start viewagent upgrade on ubuntu 14.04 with upstart"
author "Aurora"

start on (started networking)

script
  sleep 30
  exec /usr/sbin/install_viewagent
  sleep 5
  reboot
end script