# deb/vmware-horizonagent-linux-x86_64-7.1.0/opt/vmware/VMware-horizonagent-linux-x86_64-7.1.0-5141896-config/upgradeHorizonAgent.service
[Unit]
Description=upgrade VMWare Horizon Agent 
After=network-online.target
After=multi-user.target

[Service]
Type=idle
User=root
PIDFile=/run/upgradeHorizonAgent.pid
ExecStart=/usr/sbin/install_viewagent

[Install]
WantedBy=multi-user.target